# build/sysroot/usr/lib/systemd/system/cassetteos-local-storage-first.service
[Unit]
Before=docker.service
Description=Hello World

[Service]
ExecStart=/usr/bin/cassetteos-local-storage -init
Type=oneshot


[Install]
WantedBy=multi-user.target